# mikrotik-bot.service
[Unit]
Description=MikroTik Telegram Bot
After=network.target

[Service]
Type=simple
User=zesbe
WorkingDirectory=/home/zesbe/mikrotik-bot
ExecStart=/usr/bin/python3 /home/zesbe/mikrotik-bot/bot.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target